# systemd/cpanel-dyndns.service
[Unit]
Description=Run cpanel-dyndns
Wants=cpanel-dyndns.timer

[Service]
ExecStart=

[Install]
WantedBy=multi-user.target